# deploy/winit-no-sales-alert.timer.example
# 中国时间（北京时间）每天 10:00 发无动销飞书（建议在当日 06:00 入库完成之后）
#
# 与 winit-daily-sync.timer 相同：OnCalendar 用主机本地时区，请先：
#   sudo timedatectl set-timezone Asia/Shanghai

[Unit]
Description=Timer for winit-no-sales-alert.service (10:00 local = Beijing if TZ is Asia/Shanghai)

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
