# deploy/winit-no-sales-alert.service.example
# 无动销定时任务：读 SQLite 最新快照 → 按规则统计 → 飞书「no_sales」Webhook
# 默认定时见 winit-no-sales-alert.timer.example：本地 10:00 + 时区 Asia/Shanghai = 北京时间 10:00
# 依赖：当日数据已由 winit-daily-sync（或手动 run_daily）入库后再跑本任务
# .env：WINIT_FEISHU_WEBHOOK_NO_SALES、WINIT_PUBLIC_BASE_URL
#
# sudo cp deploy/winit-no-sales-alert.service.example /etc/systemd/system/winit-no-sales-alert.service
# sudo cp deploy/winit-no-sales-alert.timer.example /etc/systemd/system/winit-no-sales-alert.timer
# sudo systemctl daemon-reload && sudo systemctl enable --now winit-no-sales-alert.timer

[Unit]
Description=winit no-sales morning alert (Feishu)
After=network.target

[Service]
Type=oneshot
# 飞书 11232 频率限制时默认最多额外重试 3 次、间隔 45s，总等待可能达数分钟；勿用默认 90s 超时杀掉进程
TimeoutStartSec=300
User=root
WorkingDirectory=/opt/winit-analytics
Environment="PATH=/opt/winit-analytics/.venv/bin"
StandardOutput=journal
StandardError=journal
ExecStart=/opt/winit-analytics/.venv/bin/python /opt/winit-analytics/run_no_sales_morning_job.py

[Install]
WantedBy=multi-user.target
